# targets/silphcoanalytics/agent-fleet-watch.service
[Unit]
Description=SilphCo Agent Fleet watcher (issue dispatch + PR loop)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/evan/Documents/silphcoanalytics
EnvironmentFile=/home/evan/Documents/silphcoanalytics/.env
ExecStart=/home/evan/.local/bin/agent-fleet-watch --workspace /home/evan/Documents/silphcoanalytics
Restart=always
RestartSec=5
# KillMode=mixed: SIGTERM only the main watcher (graceful), then SIGKILL the
# entire cgroup after TimeoutStopSec. Previous KillMode=process leaked
# cursor-sdk-bridge / fastmcp children across restarts (they survived as
# orphans but the bridge's parent was dead, so they sat idle forever).
KillMode=mixed
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-fleet-watch

[Install]
WantedBy=default.target
