# examples/agent-fleet-pr-loop.service
[Unit]
Description=Agent Fleet PR loop watcher (review-fix-merge)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/absolute/path/to/your/repo
Environment=CURSOR_API_KEY=%i
ExecStart=/usr/bin/agent-fleet-pr-loop --workspace /absolute/path/to/your/repo
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
